# templates/burp-ui.service.erb
[Unit]
Description=burp-ui init script
After=network.target

[Service]
User=<%= scope.lookupvar('burp::ui::user') %>
Group=<%= scope.lookupvar('burp::ui::group') %>
ExecStart=<%= @binary_path %> --config <%= scope.lookupvar('burp::ui::config_file')%> --logfile /var/log/burp/ui.log -vvvv --debug
Type=idle

[Install]
WantedBy=multi-user.target
