# templates/burp-ui.service.erb
[Unit]
Description=<%= @service_name %> init script
After=network.target

[Service]
User=<%= scope.lookupvar('burp::ui::user') %>
Group=<%= scope.lookupvar('burp::ui::group') %>
ExecStart=<%= @binary_path %> <%= @daemon_args %>
Type=idle

[Install]
WantedBy=multi-user.target
